# deploy/kiosk/meshsat-ble-addr-pin.service
[Unit]
Description=Pin Pi 5 BT adapter flags so Reticulum GATT advertises on the public BD_ADDR (MESHSAT-678)
# host-ops: allowed-in-standalone — must run on host, not in meshsat container.
# The meshsat container's AF_BLUETOOTH HCI_CHANNEL_CONTROL socket returns
# an empty controller list even with privileged:true + network_mode:host,
# so btmgmt mutations have to happen from the host's own mgmt socket.
#
# After=bluetooth.service: Running pre-bluetoothd on Pi 5 brcmfmac blocks
# indefinitely on `btmgmt power off` because the kernel mgmt socket
# isn't responsive before bluetoothd drives the controller through its
# init path (field-verified stuck 6+ hours on both kits 2026-04-24).
# Post-bluetoothd the adapter is already up + responsive; we only flip
# flags. The script is idempotent and wrapped in per-step timeouts.
After=bluetooth.service
Wants=bluetooth.service
Documentation=https://gitlab.nuclearlighters.net/products/cubeos/meshsat MESHSAT-678

[Service]
Type=oneshot
RemainAfterExit=yes
# Hard timeout cap — the script's own timeout 8 per-step plus the
# readiness poll at the top can't exceed ~90 s. Anything longer is a
# stuck mgmt socket we don't want blocking graphical.target.
TimeoutStartSec=90
ExecStart=/usr/local/bin/meshsat-ble-addr-pin

[Install]
WantedBy=multi-user.target
