# scripts/99-meshsat-aioc.rules
# /etc/udev/rules.d/99-meshsat-aioc.rules
# AIOC (All-In-One Cable) hot-plug handler for MeshSat field kits [MESHSAT-403]
#
# When AIOC is plugged in (or re-enumerates after USB reset):
#   1. Set hidraw permissions for CM108 PTT (Direwolf needs access)
#   2. Restart Direwolf to pick up the audio device
#
# When AIOC is unplugged:
#   Direwolf will fail on its own; systemd Restart=always handles it

# AIOC USB Audio interface — restart Direwolf on re-plug
ACTION=="add", SUBSYSTEM=="sound", ATTRS{idVendor}=="1209", ATTRS{idProduct}=="7388", \
  TAG+="systemd", RUN+="/bin/systemctl restart direwolf.service"

# AIOC HID interface (CM108 PTT) — set permissions for Direwolf
ACTION=="add", SUBSYSTEM=="hidraw", ATTRS{idVendor}=="1209", ATTRS{idProduct}=="7388", \
  MODE="0660", GROUP="audio"
