# scripts/direwolf.service
# /etc/systemd/system/direwolf.service
# Direwolf APRS TNC — self-healing service for MeshSat field kits
# [MESHSAT-403]
#
# Self-healing layers:
#   1. ExecStartPre: pre-flight ALSA test + USB reset if needed
#   2. Restart=always: recovers from crashes, segfaults, audio errors
#   3. BindsTo=sys-devices-*-AIOC: auto-stop when USB unplugged
#   4. udev rule (99-meshsat-aioc.rules): auto-start when USB plugged

[Unit]
Description=Direwolf APRS TNC (KISS TCP for MeshSat)
After=sound.target
# Direwolf needs the ALSA sound system ready
Wants=sound.target

[Service]
Type=simple

# Pre-flight: test ALSA, reset USB if stuck, kill stale holders
ExecStartPre=/bin/bash /usr/local/bin/direwolf-preflight.sh

ExecStart=/usr/bin/direwolf -t 0 -c /etc/direwolf.conf

# Self-healing: restart on ANY exit (crash, error, signal, clean)
Restart=always
RestartSec=5

# Limit restart burst: max 5 restarts in 60s, then wait 120s
StartLimitIntervalSec=120
StartLimitBurst=5

# Systemd watchdog: if Direwolf doesn't notify within 120s, restart
# Note: Direwolf doesn't support sd_notify — the WatchdogSec only
# triggers if the process exits. For stuck-process detection, the
# bridge-level silence watchdog (Layer 3) handles it.

# Resource limits
MemoryMax=64M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=direwolf

[Install]
WantedBy=multi-user.target
